# src/service/garage/run
#!/bin/sh

cd /home/pi/garage/garagemonitor
exec python app.py 2>&1
